# python/ansible/hadoop_elk/roles/flume/templates/init/flume-syslog.conf
# testservice - test service job file

description "Flume Syslog Agent using MemChannel and HDFS"
author "Rusty Hann <rustyhann@gmail.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas

# When to start the service
start on (started hadoop-dfs and started hadoop-yarn and runlevel [2345])

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Run the service as another user
setuid hadoop

# Specify working directory
chdir /usr/local/flume

# Specify the process/command to start, e.g.
script

    logger -t FLUME-SYSLOG "Loading /etc/.profile"
    . /etc/profile

    logger -t FLUME-SYSLOG "Starting Flume with syslog_mem_hdfs agent"
    exec bin/flume-ng agent --conf ./conf/ -f conf/agent-syslog_mem_hdfs.conf -n syslog_mem_hdfs

end script
